# systemd/claude-inbox.timer
[Unit]
Description=Poll the Claude job inbox (CIFS-safe trigger)

[Timer]
OnBootSec=30s
OnUnitActiveSec=15s
AccuracySec=2s

[Install]
WantedBy=timers.target
