# systemd/claude-heartbeat.timer
[Unit]
Description=Claude worker heartbeat (every 30s)

[Timer]
OnBootSec=10s
OnUnitActiveSec=30s
AccuracySec=2s

[Install]
WantedBy=timers.target
