# gitops/bullpen-gitops.service
[Unit]
Description=Bullpen gitops — pull origin/main and redeploy this worker
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bullpen-gitops
